# supervisord-heartbeat.conf
[program:heartbeat]
command=/data/heartbeat.sh
numprocs=1
autostart=true
autorestart=true
