# etc/supervisord-apache2.conf
[program:apache2]
command=/data/start-apache2.sh
numprocs=1
autostart=true
autorestart=true
